# src/etc/s6/tor/run
#!/bin/sh
set -e

DATADIR=/var/lib/tor

chmod -R 0600 $DATADIR
chmod 0700 $DATADIR

chown -R nobody:nogroup $DATADIR

CONFIG=/etc/tor/torrc

touch $CONFIG
chmod 0644 $CONFIG

printenv | grep '^TOR_' | sed -E 's/TOR_([^=]*)=(.*)/\1 \2/' > $CONFIG
cat << EOF >> $CONFIG
SOCKSPort 127.0.0.1:9050
DataDirectory $DATADIR
EOF

exec s6-setuidgid nobody tor -f $CONFIG
